# etc/systemd/system/satechi.service
[Unit]
Description=Satechi Remote service
BindsTo=dev-SBMMR.device

[Service]
ExecStart=/usr/local/bin/satechi.py
